# src/asynserver.service
[Unit]
Description=A simple asynchronous server using epoll
After=network.target

[Service]
Type=simple
ExecStart=/usr/local/bin/asynserver
Restart=on-failure

[Install]
WantedBy=multi-user.target
